# Module1/crontab
00 3    * * *  root    bash /bin/backup_script.sh